# websocket-server.service
[Unit]
Description=WebSocket Server Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/live.onlineiotdata.in
ExecStart=/usr/bin/php websocket_server.php
Restart=always
RestartSec=10
StandardOutput=append:/var/www/live.onlineiotdata.in/logs/websocket.log
StandardError=append:/var/www/live.onlineiotdata.in/logs/websocket.log

[Install]
WantedBy=multi-user.target 